# deploy/susanino-bot.service
[Unit]
Description=Susanino Map Telegram Bot
After=network.target

[Service]
Type=simple
User=susanino
WorkingDirectory=/opt/susanino-map/bot
EnvironmentFile=/opt/susanino-map/bot/.env
ExecStart=/opt/susanino-map/venv/bin/python bot.py
Restart=always
RestartSec=5
# Базовый харденинг
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
